# deploy/unifiedops-listener@.service
; ============================================================================
; UnifiedOps v2 Listener Template Service
;
; Responsibilities:
;   - Runs a specific syslog trap listener script natively via Python.
;   - Connects to the local or remote InfluxDB defined in its .env file.
;   - Injects the %i specifier to dynamically pick the Python script and 
;     its corresponding .env configuration file.
;
; Install:
;   sudo cp deploy/unifiedops-listener@.service /etc/systemd/system/
;   sudo systemctl daemon-reload
;
; Usage:
;   Assuming you have a listener script named `syslog_trap_listener_bcp.py`:
;   1. Create its config file: `/etc/unifiedops/listener-syslog_trap_listener_bcp.env`
;   2. Enable and start it:
;      `sudo systemctl enable --now unifiedops-listener@syslog_trap_listener_bcp`
;
; ============================================================================
[Unit]
Description=UnifiedOps v2 Listener (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unifiedops
Group=unifiedops
WorkingDirectory=/opt/unifiedops
EnvironmentFile=/etc/unifiedops/listener-%i.env
ExecStart=/opt/unifiedops/.venv/bin/python /opt/unifiedops/listener/%i.py
Restart=on-failure
RestartSec=4
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/unifiedops

[Install]
WantedBy=multi-user.target
